# configuration_systemd/frases.timer
[Unit]
Description=Troca frase diariamente às 07h

[Timer]
OnCalendar=07:00
Persistent=true

[Install]
WantedBy=timers.target
